# unix-remote/org.genera.lpdd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>            <string>org.genera.lpdd</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/ldbeth/Public/Projects/linux-vlm/unix-remote/lpdd</string>
    <string>--spool</string><string>/Users/ldbeth/genera-spool</string>
    <!-- bind the vmnet bridge IP (host side of og2vlm/.VLM), never 0.0.0.0;
         --wait retries the bind until that address exists, so the job rides
         the VLM coming up and going down instead of throttle-failing. -->
    <string>--host</string><string>192.168.2.1</string>
    <string>--wait</string>
    <string>-v</string>
  </array>
  <!-- A long-running server (NOT inetdCompatibility): launchd keeps it alive,
       lpdd itself handles the port-515 listen and waits for the bridge. -->
  <key>KeepAlive</key><true/>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/ldbeth/genera-spool/lpdd.log</string>
  <key>StandardErrorPath</key><string>/Users/ldbeth/genera-spool/lpdd.log</string>
</dict>
</plist>
